# rootfs/s6/.crond/run
#!/bin/sh

set -euo pipefail

exec crond -f -d "${CROND_LOG_LEVEL}" ${CROND_EXTRA_ARGS:-}
