# rootfs/s6/unbound/run
#!/bin/sh

set -euo pipefail

exec unbound -d ${UNBOUND_EXTRA_ARGS:-}
